# com.user.signal-claude-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.signal-claude-bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/rossbrown/Desktop/DevProjects/Claude/claude-signal-integration/.venv/bin/python3</string>
        <string>/Users/rossbrown/Desktop/DevProjects/Claude/claude-signal-integration/src/daemon/service.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/rossbrown/Desktop/DevProjects/Claude/claude-signal-integration</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/rossbrown/Library/Logs/signal-claude-bot/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/rossbrown/Library/Logs/signal-claude-bot/stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/rossbrown/Desktop/DevProjects/Claude/claude-signal-integration</string>
    </dict>
</dict>
</plist>
